# mainsail.service
#Systemd service file for klipper instance banana
#This file has been automatically generated by klipper-deployer. Please do not edit it manually.

[Unit]
Description=Starts klipper on startup
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=migration
RemainAfterExit=yes
ExecStart=/home/migration/klipper-deployer/klippy-env/bin/python /home/migration/klipper-deployer/repos/klipper/klippy/klippy.py /home/migration/klipper-deployer/instances/banana/printer_data/config/printer.cfg -l /home/migration/klipper-deployer/instances/banana/printer_data/logs/klippy.log -a /tmp/klippy_uds_undefined_banana
Restart=always
RestartSec=10
